# sovereign-swarm/systemd/pelican-power.timer
[Unit]
Description=Run Pelican Power Manager periodically
Documentation=https://github.com/Strategickhaos/sovereign-swarm

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
