# sovereign-swarm/systemd/pelican-power.service
[Unit]
Description=Pelican Power Manager
Documentation=https://github.com/Strategickhaos/sovereign-swarm
After=network.target

[Service]
Type=oneshot
ExecStart=/opt/sovereign-swarm/power-manager.sh optimize

[Install]
WantedBy=multi-user.target
